# install/DEBIAN/prerm
#!/bin/sh
systemctl disable picp
systemctl stop picp
exit 0